# kong.service
[Unit]
Description= kong service
After=syslog.target network.target

[Service]
User=root
Group=root
Type=forking
ExecStart=/usr/local/bin/kong start -c /usr/local/kong/conf/kong.conf
ExecReload=/usr/local/bin/kong reload -c /usr/local/kong/conf/kong.conf
ExecStop=/usr/local/bin/kong stop

[Install]
WantedBy=multi-user.target
